# systemd/scrz-node-supervisor.service
[Unit]
Description=scrz node supervisor

[Service]
Environment=LD_LIBRARY_PATH=/opt/scrz/lib
ExecStart=/opt/scrz/bin/scrz-node-supervisor
